# roles/prometheus/templates/prometheus_consul_exporter_supervisord.conf.j2
[program:prometheus_consul_exporter]
command= /opt/bin/start_prometheus_consul_exporter.sh
autostart=true
autorestart=true
startretries=999
redirect_stderr=True
nodaemon = False
debug = False
stdout_logfile = /var/tmp/prometheus_consul_exporter.out
stderr_logfile = /var/tmp/prometheus_consul_exporter.err
loglevel = "info"
stopasgroup=true
